# deployment/Linux/Debian/debian/DEBIAN/postinst
#!/bin/bash
systemctl enable ideservice.service
systemctl start ideservice.service

